# scripts/macos-guest-agent/container-macos-guest-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apple.container.macos.guest-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>
        <string>--port</string>
        <string>__PORT__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/container-macos-guest-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/container-macos-guest-agent.log</string>
</dict>
</plist>
